# deployment/track2college.timer
[Unit]
Description=Run Track2College Pipeline Monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
