# deployment/track2college.service.template
[Unit]
Description=Track2College Pipeline Automation

[Service]
Type=oneshot
WorkingDirectory={{PROJECT_DIR}}
ExecStart={{PYTHON_EXEC}} run_pipeline.py
StandardOutput=journal
StandardError=journal
